# ops/etc/systemd/system/dandorak-post-reboot.service
[Unit]
Description=Dandorak post-reboot Auto Deploy restore
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/dandorak-post-reboot.sh
RemainAfterExit=no
KillMode=process

[Install]
WantedBy=multi-user.target
